# github-treator.socket
[Socket]
ListenStream=/run/github-commands.sk

[Install]
WantedBy=sockets.target