# github-treator.service
[Service]
Type=oneshot
ExecStart=/bin/node /home/lenny/gits/machine-github-treator/index.js
StandardOutput=syslog
User=lenny
Group=users